# dev/soultribe-gunicorn.service
# /etc/systemd/system/soultribe-gunicorn@.service
# Adjust WorkingDirectory and paths to .venv as needed.
[Unit]
Description=SoulTribe FastAPI (gunicorn + uvicorn workers)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Instance-based: pass your Linux username as the instance, e.g. soultribe-gunicorn@YOURUSER
User=%i
Group=%i
WorkingDirectory=/home/%i/windsurf/soultribe.chat
EnvironmentFile=/home/%i/windsurf/soultribe.chat/.env

# Gunicorn with Uvicorn workers
ExecStart=/home/%i/windsurf/soultribe.chat/.venv/bin/gunicorn \
  -k uvicorn.workers.UvicornWorker \
  --workers 2 \
  --threads 2 \
  --bind 0.0.0.0:8001 \
  --timeout 60 \
  --access-logfile - \
  --error-logfile - \
  main:app

Restart=always
RestartSec=3
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
