# dev/soultribe-gunicorn.user.service
# ~/.config/systemd/user/soultribe-gunicorn.service
# User-level systemd unit (no sudo). Uses %h for your home directory.
[Unit]
Description=SoulTribe FastAPI (gunicorn + uvicorn workers, user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/windsurf/soultribe.chat
EnvironmentFile=%h/windsurf/soultribe.chat/.env
ExecStart=%h/windsurf/soultribe.chat/.venv/bin/gunicorn \
  -k uvicorn.workers.UvicornWorker \
  --workers 2 \
  --threads 2 \
  --bind 0.0.0.0:8001 \
  --timeout 60 \
  --access-logfile - \
  --error-logfile - \
  main:app
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
